# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm -c /etc/

[program:nginx]
command=/bin/bash -c nginx
